# omd/packages/check_mk/skel/etc/cron.d/cmk_transform_inventory_trees
# Once per day at 00:00
0 0 * * * cmk-transform-inventory-files
