# omd/packages/check_mk/skel/etc/cron.d/cmk_collect_transmit_product_usage_analytics
# Run every 30 minutes
# Note: in store_next_run_ts() the next (random) time point is set between 30 and 31 days from now
*/30 * * * * ###ROOT###/bin/cmk-product-usage --cron
