# omd/packages/check_mk/skel/etc/cron.d/cmk_update_license_usage
# Run every 30 minutes between 8am and 5pm
# Note: in update_license_usage() the next (random) time point is set between 8am and 4pm
*/30 8-17 * * * ###ROOT###/bin/cmk-update-license-usage
