# omd/packages/check_mk/skel/etc/cron.d/cmk_discovery
# Every 5 minutes autodiscover services marked by discovery check
*/5 * * * * cmk --discover-marked-hosts
