# omd/packages/maintenance/skel/etc/cron.d/diskspace
#
# Run every hour, 5 minutes after full hour
#
5 * * * * diskspace > $OMD_ROOT/var/log/diskspace.log 2>&1
