# agents/systemd/check_mk@.service
[Unit]
Description=Check_MK

[Service]
ExecStart=/usr/bin/check_mk_agent

User=root
Group=root

StandardInput=socket
