# omd/packages/maintenance/skel/etc/init.d/crontab
#!/bin/bash
# 
# Start only if Crontab is enabled
. ###ROOT###/etc/omd/site.conf

NAME="crontab"
CROND="###ROOT###/etc/cron.d/*"
CRONTAB_OPTS=""
USER="###SITE###"

# check for root
# Workaround for http://omdistro.org/issues/157
if [ `id -u` -eq "0" ]; then
	CRONTAB_OPTS="-u ###SITE###"
fi

_merge_crontabs() {
    cat<< HERE
#
# Do not edit this file. It will be recreated each time OMD
# is started or reloaded.
#
# execute 'omd reload crontab'
# to rebuild this file out of ${OMD_ROOT}/etc/cron.d/*
#
# --ENVIRONMENT------------------------------------------------
SHELL=/bin/bash
BASH_ENV=${OMD_ROOT}/.profile
OMD_ROOT=${OMD_ROOT}
OMD_SITE=${OMD_SITE}
PATH=${PATH}
MAILTO="${CONFIG_ADMIN_MAIL}"
# ------------------------------------------------------------
HERE
    cat ${CROND}
}

# See how we were called.
case "$1" in

    start)
	echo -en "Starting crontab..."
	_merge_crontabs ${CROND} | crontab $CRONTAB_OPTS - > /dev/null
        if [ $? -eq 0 ]; then
            echo "OK"
            exit 0
	else
	    echo "failed"
	    exit 1
        fi
        ;;

    stop)
        echo -n "Stopping crontab..."
        crontab $CRONTAB_OPTS -r || exit 0
        echo OK
        exit 0
        ;;

    restart|reload)
        $0 stop
        $0 start
        ;;
    status)
        crontab $CRONTAB_OPTS -l >/dev/null 2>&1
        if [ $? -eq 0 ]; then
            echo "crontab initialized"
            exit 0
        else
            echo "crontab is empty"
            exit 1
        fi
        ;;
    merge)
        _merge_crontabs ${CROND}
        exit 0
        ;;
    *)
        echo "Usage: crontab {start|stop|restart|reload|status}"
        exit 1
        ;;

esac
  
# End of this script
